# supervisord.conf
; supervisor config file

[unix_http_server]
file=/var/run//supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[inet_http_server]         ; inet (TCP) server disabled by default
port=0.0.0.0:9001        ; (ip_address:port specifier, *:port for all iface)
username=admin              ; (default is no username (open server))
password=Admin123$           ; (default is no password (open server))

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run//supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.


[program:jsterm]
command= python -m SimpleHTTPServer 80
directory=/home/jsterm/ ; directory to cwd to before exec (def no cwd)
autostart=true ; start at supervisord start (default: true)
autorestart=true ; whether/when to restart (default: unexpected)
startsecs=3 ; number of secs prog must stay running (def. 1)
redirect_stderr=true ; redirect proc stderr to stdout (default false) 
stdout_logfile=/tmp/jsterm.log ; stdout log path, NONE for none; default AUTO, log out
stopasgroup=true

[program:v2ex_django]
command= python app.py
directory=/home/v2exfriends/ ; directory to cwd to before exec (def no cwd)
autostart=true ; start at supervisord start (default: true)
autorestart=true ; whether/when to restart (default: unexpected)
startsecs=3 ; number of secs prog must stay running (def. 1)
redirect_stderr=true ; redirect proc stderr to stdout (default false) 
stdout_logfile=/tmp/v2ex_django.log ; stdout log path, NONE for none; default AUTO, log out
stopasgroup=true
;; environment = NEW_RELIC_CONFIG_FILE="/home/v2exfriends/newrelic.ini"

[program:v2ex_django_celery_beat]
command= celery -A tasks beat
directory=/home/v2exfriends/ ; directory to cwd to before exec (def no cwd)
autostart=true ; start at supervisord start (default: true)
autorestart=true ; whether/when to restart (default: unexpected)
startsecs=3 ; number of secs prog must stay running (def. 1)
redirect_stderr=true ; redirect proc stderr to stdout (default false) 
stdout_logfile=/tmp/v2ex_django_celery_beat.log ; stdout log path, NONE for none; default AUTO, log out
stopasgroup=true
;; environment = NEW_RELIC_CONFIG_FILE="/home/v2exfriends/newrelic.ini"

[program:v2ex_django_celery]
command= celery -A tasks worker -l info -c 400 -P gevent
directory=/home/v2exfriends/ ; directory to cwd to before exec (def no cwd)
autostart=true ; start at supervisord start (default: true)
autorestart=true ; whether/when to restart (default: unexpected)
startsecs=3 ; number of secs prog must stay running (def. 1)
redirect_stderr=true ; redirect proc stderr to stdout (default false) 
stdout_logfile=/tmp/v2ex_django_celery.log ; stdout log path, NONE for none; default AUTO, log out
stopasgroup=true
;; environment = NEW_RELIC_CONFIG_FILE="/home/v2exfriends/newrelic.ini"

[program:v2ex_flower]
command= celery flower -A tasks --broker=redis://localhost:6379/0 --basic_auth=admin:admin
directory=/home/v2exfriends/ ; directory to cwd to before exec (def no cwd)
autostart=true ; start at supervisord start (default: true)
autorestart=true ; whether/when to restart (default: unexpected)
startsecs=3 ; number of secs prog must stay running (def. 1)
redirect_stderr=true ; redirect proc stderr to stdout (default false) 
stdout_logfile=/tmp/v2ex_flower.log ; stdout log path, NONE for none; default AUTO, log out
stopasgroup=true
;; environment = NEW_RELIC_CONFIG_FILE="/home/v2exfriends/newrelic.ini"

[include]
files = /etc/supervisor/conf.d/*.conf
